# scripts/digests/_retired_plists/com.hatchingsolutions.digest-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hatchingsolutions.digest-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kurtishatcher/intelligence-dashboard/scripts/digest-sync.sh</string>
        <string>start</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>LAUNCHED_BY_LAUNCHD</key>
        <string>1</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/kurtishatcher/intelligence-dashboard/scripts/digests/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kurtishatcher/intelligence-dashboard/scripts/digests/sync.log</string>
</dict>
</plist>
